# supervisor/snin_mesh.conf
[group:snin]
programs=route_engine,crv2,smart_router,bridge,cheque_book,verifier,gossip_0,gossip_1,gossip_2,gossip_3,gossip_4,external_gateway

[program:route_engine]
command=python3 -u /home/agent/data/sites/relay-mesh/route_engine.py
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/route_engine.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/route_engine.err
user=agent

[program:crv2]
command=python3 -u /home/agent/data/sites/relay-mesh/content_router_v2.py 9920
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/content_router.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/content_router.err
user=agent

[program:smart_router]
command=python3 -u /home/agent/data/sites/relay-mesh/smart_router.py
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/smart_router.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/smart_router.err
user=agent

[program:bridge]
command=python3 -u /home/agent/data/sites/relay-mesh/nostr_bridge.py
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/nostr_bridge.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/nostr_bridge.err
user=agent

[program:cheque_book]
command=python3 -u /home/agent/data/sites/relay-mesh/cheque_book.py
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/cheque_book.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/cheque_book.err
user=agent

[program:verifier]
command=python3 -u /home/agent/data/sites/relay-mesh/verifier.py
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/verifier.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/verifier.err
user=agent

[program:gossip_0]
command=python3 -u /home/agent/data/sites/relay-mesh/gossip_shard.py --port 9100 --shard-id 0 --agents 100
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=3
stopsignal=TERM
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_0.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_0.err
user=agent

[program:gossip_1]
command=python3 -u /home/agent/data/sites/relay-mesh/gossip_shard.py --port 9101 --shard-id 1 --agents 100
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=3
stopsignal=TERM
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_1.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_1.err
user=agent

[program:gossip_2]
command=python3 -u /home/agent/data/sites/relay-mesh/gossip_shard.py --port 9102 --shard-id 2 --agents 100
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=3
stopsignal=TERM
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_2.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_2.err
user=agent

[program:gossip_3]
command=python3 -u /home/agent/data/sites/relay-mesh/gossip_shard.py --port 9103 --shard-id 3 --agents 100
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=3
stopsignal=TERM
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_3.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_3.err
user=agent

[program:gossip_4]
command=python3 -u /home/agent/data/sites/relay-mesh/gossip_shard.py --port 9104 --shard-id 4 --agents 100
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=3
stopsignal=TERM
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_4.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/gossip_4.err
user=agent

[program:external_gateway]
command=python3 -u /home/agent/data/sites/relay-mesh/external_gateway.py
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/external_gateway.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/external_gateway.err
user=agent


[group:snin_agents]
programs=agent_forecaster,agent_archivist,agent_anton

[program:agent_forecaster]
command=python3 -u /home/agent/data/sites/relay-mesh/agent_daemon.py forecaster_ai
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/agent_forecaster.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/agent_forecaster.err
user=agent

[program:agent_archivist]
command=python3 -u /home/agent/data/sites/relay-mesh/agent_daemon.py archivist_ai
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/agent_archivist.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/agent_archivist.err
user=agent

[program:agent_anton]
command=python3 -u /home/agent/data/sites/relay-mesh/agent_daemon.py anton_ai
directory=/home/agent/data/sites/relay-mesh
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/agent/data/sites/relay-mesh/logs/agent_anton.log
stderr_logfile=/home/agent/data/sites/relay-mesh/logs/agent_anton.err
user=agent
